# templates/default/rhel-init-service.erb
#!/bin/sh
#
# processname: <%= @service_name %>
# description: <%= @service_desc %>
# chkconfig: <% if @run[:daemon] %>2345<%else %>-<% end %> 50 50

. /etc/init.d/functions

case "${1}" in
    start)
        echo_passed "Starting..."
        mkdir --mode=0770 --parents <%= @piddir %>
        chown <%= @user %>:<%= @group %> <%= @piddir %>

        <% if scalrpy2? node %>
        daemon --force --user <%= @user %> --pidfile <%= @pidfile %> <%= @executable %> -- <%= node[:scalr][:core][:location] %>/app/python/scalrpy/<%= @service_module %>.py --pid-file=<%= @pidfile %> --log-file=<%= @logfile %> --config=<%= node[:scalr][:core][:configuration] %> --verbosity=INFO --daemon <%= @daemon_extra_args %> start
        <% else %>
        daemon --force --user <%= @user %> --pidfile <%= @pidfile %> <%= @executable %> -m scalrpy.<%= @service_module %> -p <%= @pidfile %> -l <%= @logfile %> -c <%= node[:scalr][:core][:configuration] %> -vvv --start <%= @service_extra_args %>
        <% end %>
        ;;

    stop)
        <% if scalrpy2? node %>
        <%= @executable %> <%= node[:scalr][:core][:location] %>/app/python/scalrpy/<%= @service_module %>.py --pid-file=<%= @pidfile %> --log-file=<%= @logfile %> --config=<%= node[:scalr][:core][:configuration] %> stop
        <% else %>
        # Don't use the buit-in killproc function here, because we have an entire group to kill
        echo_passed "Stopping..."
        [ -f "<%= @pidfile %>" ] && kill -TERM -$(ps -o pgid= $(cat -- "<%= @pidfile %>") | grep -o [0-9]*)
        ret="$?"
        rm -f -- "<%= @pidfile %>"
        <% end %>
        exit $?
        ;;

    restart)
        ${0} stop
        sleep 1
        ${0} start
        ;;

    status)
        status -p <%= @pidfile %> <%= @executable %>  <% @service_name %>
        ;;

    *)
        echo "Usage: ${0} {start|stop|status|restart}"
        exit 1
        ;;
esac
