# templates/default/rhel-init-service.erb
#!/bin/sh
#
# processname: <%= @service_name %>
# description: <%= @service_desc %>
# chkconfig: <% if @run[:daemon] %>2345<%else %>-<% end %> 50 50

. /etc/init.d/functions

case "${1}" in
   start)
      echo_passed "Starting..."
      mkdir --mode=0770 --parents <%= @piddir %>
      chown <%= @user %>:<%= @group %> <%= @piddir %>

      daemon --force --user <%= @user %> --pidfile <%= @pidfile %> <%= @executable %> -- <%= node[:scalr][:core][:location] %>/app/python/scalrpy/<%= @service_script %> --pid-file=<%= @pidfile %> --log-file=<%= @logfile %> --config=<%= node[:scalr][:core][:configuration] %> --verbosity=3 --daemon <%= @daemon_extra_args %> start
      ;;

   stop)
      echo_passed "Stopping..."
      <%= @executable %> <%= node[:scalr][:core][:location] %>/app/python/scalrpy/<%= @service_script %> --pid-file=<%= @pidfile %> --log-file=<%= @logfile %> --config=<%= node[:scalr][:core][:configuration] %> stop
      ;;

   restart)
      ${0} stop
      sleep 1
      ${0} start
      ;;

  status)
      status -p <%= @pidfile %> <%= @executable %>  <% @service_name %>
      ;;

   *)
      echo "Usage: ${0} {start|stop|status|restart}"
      exit 1
      ;;
esac
