# templates/default/rhel-init-service.erb
#!/bin/sh
#
# processname: <%= @daemon_name %>
# description: <%= @daemon_desc %>
# chkconfig: <% if @run[:daemon] %>2345<% end %> 50 50

. /etc/init.d/functions

case "${1}" in
   start)
      echo_passed "Starting..."
      mkdir --mode=0770 --parents <%= @piddir %>
      chown <%= @user %>:<%= @group %> <%= @piddir %>

      daemon --force --user <%= @user %> --pidfile <%= @pidfile %> <%= @executable %> -m scalrpy.<%= @daemon_module %> -p <%= @pidfile %> -l <%= @logfile %> -c <%= node[:scalr][:core][:configuration] %> -vvv --start <%= @daemon_extra_args %>
      ;;

   stop)
      echo_passed "Stopping..."
      killproc -p <%= @pidfile %> <%= @executable %> SIGINT
      ;;

   restart)
      ${0} stop
      sleep 1
      ${0} start
      ;;

   *)
      echo "Usage: ${0} {start|stop|restart}"
      exit 1
      ;;
esac

# End scriptname
