# templates/default/upstart-service.conf.erb
description "<%= @daemon_desc %>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

expect daemon

console none

pre-start script
  if [ ! -r <%= node[:scalr][:core][:configuration] %> ]; then
    logger -is -t "\$UPSTART_JOB" "ERROR: Config file is not readable"
    exit 1
  fi
  mkdir --mode=770 --parents <%= node[:scalr][:core][:pid_dir] %> <%= node[:scalr][:core][:log_dir] %>
  chown <%= node[:scalr][:core][:users][:service] %>:<%= node[:scalr][:core][:group] %> <%= node[:scalr][:core][:pid_dir] %> <%= node[:scalr][:core][:log_dir] %>
end script

exec start-stop-daemon --start --chuid <%= node[:scalr][:core][:users][:service] %>:<%= node[:scalr][:core][:group] %> --pidfile <%= node[:scalr][:core][:pid_dir] %>/<%= @daemon_name %>.pid --exec <%= @python %> -- -m scalrpy.<%= @daemon_module %> -p <%= node[:scalr][:core][:pid_dir] %>/<%= @daemon_name %>.pid -l <%= node[:scalr][:core][:log_dir] %>/<%= @daemon_name %>.log -c <%= node[:scalr][:core][:configuration] %> -vvv --start <%= @daemon_extra_args %>
