# templates/ubuntu/init-service.erb
#!/bin/sh
########################################################################
# Begin <%= @daemon_name %>
#
# Description : <%= @daemon_desc %>
#
########################################################################

### BEGIN INIT INFO
# Provides:            <%= @daemon_name %>
# Required-Start: $local_fs $remote_fs
# Required-Stop: $local_fs $remote_fs
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: <%= @daemon_desc %>
### END INIT INFO

. /lib/lsb/init-functions

case "${1}" in
   start)
      log_action_begin_msg "Starting..."
      /sbin/start-stop-daemon --start --chdir "$PWD" --oknodo --chuid <%= @user %>:<%= @group %> --pidfile <%= @pidfile %> --exec <%= @executable %> -- -m scalrpy.<%= @daemon_module %> -p <%= @pidfile %> -l <%= @logfile %> -c <%= node[:scalr][:core][:configuration] %> -vvv --start <%= @daemon_extra_args %>
      ;;

   stop)
      log_action_begin_msg "Stopping..."
      killproc -p <%= @pidfile %> <%= @executable %>
      ;;

   restart)
      ${0} stop
      sleep 1
      ${0} start
      ;;

   *)
      echo "Usage: ${0} {start|stop|restart}"
      exit 1
      ;;
esac

exit 0

# End scriptname
