# files/scalr-server-cookbooks/scalr-server/templates/default/supervisor/rhel/supervisor.init.erb
#!/bin/bash
#
# /etc/rc.d/init.d/scalr
#
# Manages the scalr supervisor daemon
#
# chkconfig: 2345 90 10
# description: Supervisor is a process control utility, and is used to manage the Scalr processes.
# processname: scalr
# config: <%= node['supervisor']['conffile'] %>
# pidfile: <%= @run_dir %>/supervisord.pid

source /etc/rc.d/init.d/functions

export PATH="<%= scalr_exec_path node %>"
export TZ="<%= node[:rackspace_timezone][:config][:tz] %>"

start() {
    echo -n "Starting supervisor: "
    # status() returns 0 if and only if the service is properly started.
    status -p <%= @run_dir %>/supervisord.pid scalr > /dev/null && echo "already running" && return 0
    daemon "<%= @supervisord %> -c <%= node['supervisor']['conffile'] %> <%= @node['supervisor']['daemon_options'] %>"
    RETVAL=$?
    return $RETVAL
}

stop() {
    echo -n "Shutting down scalr: "
    status -p <%= @run_dir %>/supervisord.pid scalr > /dev/null
    # status() returns 3 if and only if the service is properly stopped.
    [ $? -eq 3 ] && echo "already stopped" && return 0
    killproc -p <%= @run_dir %>/supervisord.pid scalr
    RETVAL=$?
    return $RETVAL
}

case "$1" in
    start)
        start
    ;;
    stop)
        stop
    ;;
    status)
        status -p <%= @run_dir %>/supervisord.pid scalr
    ;;
    restart)
        # Exit with an error if the shutdown was not successful.
        stop || exit 1
        start
    ;;
    *)
        echo "Usage: <servicename> {start|stop|status|restart}"
        exit 3
    ;;
esac
